# componentes/easyengine/supervisord.conf
[program:easyengine]
command=/startStack.sh
